# deploy/systemd/vaultkeeper-update.timer
[Unit]
Description=VaultKeeper Auto-Update Timer
Requires=vaultkeeper-update.service

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
